# bot.conf
[program:cnc]
directory=/cnc
command=python3 app.py

autostart=true
autorestart=true

environment="/cnc", USER="cncbot"

stderr_logfile=/cnc/logfile_err.log
stdout_logfile=/cnc/logfile.log